# systemd/pico-scheduler.service
[Unit]
Description=Pico Scheduler Service
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/pico
ExecStart=/home/pi/.local/bin/uv run python -m services.scheduler
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
