# systemd/pico-sms.service
[Unit]
Description=Pico SMS Service
After=network.target

[Service]
User=pi
Group=pi
SupplementaryGroups=gammu
WorkingDirectory=/home/pi/Pico
ExecStart=/home/pi/.local/bin/uv run python -m services.sms
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
